# printwatch.service
[Unit]
Description=PrintWatch AI
After=multi-user.target

[Service]
WorkingDirectory=/home/pi/
User=pi
Type=simple
Restart=always
ExecStart=/usr/bin/python3 /home/pi/printwatch/main.py

[Install]
WantedBy=multi-user.target
